# script_systemctl/central-writer.service
[Unit]
Description=Laravel Central Writer
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
User=bruno
Group=bruno
WorkingDirectory=/var/www/Sitemas
ExecStart=/var/www/Sitemas/script_systemctl/run-central-writer.sh
Environment=APP_ROOT=/var/www/Sitemas
Restart=on-failure
RestartSec=5s
StandardOutput=journal
StandardError=journal

[Install]
WantedBy=multi-user.target
